# systemd-dropins/axentx-global-trend-harvester-daemon.service.d/override.conf
[Service]
Environment=TREND_CYCLE_GAP_SEC=1800
Environment=TREND_REQ_GAP_SEC=8.0
