# systemd-dropins/axentx-trend-arbitrage-daemon.service.d/override.conf
[Service]
Environment=TREND_ARB_POLL_SEC=15
Environment=BUDGET_TREND_ARBITRAGE=1500
