# systemd-dropins/axentx-github-deep-stream-daemon.service.d/resilience.conf
[Service]
Restart=always
RestartSec=30
StartLimitIntervalSec=600
StartLimitBurst=10
# Treat Supabase upstream noise as non-fatal
Environment=SUPABASE_503_NONFATAL=1
