# systemd/axentx-medium-crawler-daemon.service
[Unit]
Description=axentx medium-crawler — pulls Medium articles by tag, extracts monetization patterns, emits ideas to validator-queue
After=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=MEDIUM_POLL_SEC=1800
Environment=MEDIUM_TAGS=startup,saas,solopreneur,indie-hackers,business,startup-ideas,micro-saas,side-project,no-code,build-in-public,entrepreneurship,product-management,growth,bootstrapped
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-medium-crawler-daemon.py
Restart=always
RestartSec=60
MemoryMax=256M
TasksMax=8
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
