# systemd/axentx-usability-tester-daemon.service
[Unit]
Description=axentx usability-tester — runs Nielsen-10 + Thai-readability tests on each product
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=AXENTX_ROOT=/opt/axentx
Environment=USABILITY_TESTER_POLL_SEC=90
Environment=UX_QUEUE_TARGET=5
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-usability-tester-daemon.py
Restart=always
RestartSec=120
MemoryMax=192M
TasksMax=16
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
