# systemd/axentx-social-pain-stream-daemon.service
[Unit]
Description=axentx social-pain stream — Mastodon+Bluesky+GH-Discussions
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-social-pain-stream-daemon.py
Restart=always
RestartSec=15
MemoryMax=128M
CPUQuota=15%

[Install]
WantedBy=multi-user.target
