# systemd/axentx-agent-builder-daemon.service
[Unit]
Description=axentx agent-builder — auto-creates new daemons from agent-spec entries in shared_knowledge
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=AGENT_BUILDER_POLL_SEC=300
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-agent-builder-daemon.py
Restart=always
RestartSec=120
MemoryMax=256M
TasksMax=8
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
