# systemd/axentx-prompt-optimizer-daemon.service
[Unit]
Description=axentx prompt-optimizer — DSPy/GEPA-inspired prompt evolution from review-pass rate
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=PROMPT_OPTIMIZER_POLL_SEC=43200
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-prompt-optimizer-daemon.py
Restart=always
RestartSec=300
MemoryMax=192M
TasksMax=16
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
