# systemd/axentx-social-listener-daemon.service
[Unit]
Description=axentx social-listener — agentic tool-using pain crawler (spider-bot for blue-ocean opportunities)
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=LISTENER_POLL_SEC=480
Environment=LISTENER_MAX_HOPS=6
Environment=LISTENER_MAX_PAYLOAD=10
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-social-listener-daemon.py
Restart=always
RestartSec=60
MemoryMax=384M
TasksMax=24
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
