# systemd/axentx-hf-rag-loader-daemon.service
[Unit]
Description=axentx hf-rag-loader — pulls patterns from HF dataset → shared_knowledge (closed-loop bidir sync)
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=HF_RAG_LOADER_POLL_SEC=1800
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-hf-rag-loader-daemon.py
Restart=always
RestartSec=300
MemoryMax=768M
TasksMax=16
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
