# systemd/axentx-indiehackers-stream-daemon.service
[Unit]
Description=axentx indiehackers-stream — pulls product/pain signals into validator-queue
After=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=IH_POLL_SEC=1800
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-indiehackers-stream-daemon.py
Restart=always
RestartSec=60
MemoryMax=256M
TasksMax=8
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
