# systemd/axentx-revenue-verified-stream-daemon.service
[Unit]
Description=axentx revenue-verified-stream — emits high-signal ideas to validator-queue
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=REVENUE_STREAM_POLL_SEC=1800
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-revenue-verified-stream-daemon.py
Restart=always
RestartSec=120
MemoryMax=512M
TasksMax=8
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
