# systemd/axentx-product-synthesizer-daemon.service
[Unit]
Description=axentx product-synthesizer — proposes new products (OS+paid, Thai+global lens)
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=PRODUCT_SYNTH_POLL_SEC=1800
Environment=PRODUCT_SYNTH_MAX_PER_CYCLE=3
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-product-synthesizer-daemon.py
Restart=always
RestartSec=120
MemoryMax=256M
TasksMax=16
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
