# systemd/axentx-codebase-sync-daemon.service
[Unit]
Description=axentx codebase-sync — git pull every /opt/axentx repo every 10min (cross-VM sync)
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=AXENTX_ROOT=/opt/axentx
Environment=CODEBASE_SYNC_POLL_SEC=600
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-codebase-sync-daemon.py
Restart=always
RestartSec=20
MemoryMax=192M
TasksMax=8
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
