# systemd/axentx-llm-health-watchdog-daemon.service
[Unit]
Description=axentx llm-health-watchdog — periodic provider probe + degradation alerts
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=LLM_HEALTH_POLL_SEC=180
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-llm-health-watchdog-daemon.py
Restart=always
RestartSec=120
MemoryMax=192M
TasksMax=8
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
