# systemd/axentx-demand-amplifier-daemon.service
[Unit]
Description=axentx demand-amplifier — pull-based work cascade (no agent idle, ever)
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=AMPLIFIER_POLL_SEC=30
Environment=SHARED_QUEUES=/opt/surrogate-1-harvest/state/swarm-shared
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-demand-amplifier-daemon.py
Restart=always
RestartSec=30
MemoryMax=192M
TasksMax=32
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
