# systemd/axentx-github-ideabank-daemon.service
[Unit]
Description=axentx github-ideabank — emits high-signal ideas to validator-queue
After=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
Environment=IDEABANK_POLL_SEC=86400
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-github-ideabank-daemon.py
Restart=always
RestartSec=120
MemoryMax=512M
TasksMax=8
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
