# systemd/axentx-biz-pipeline-daemon.service
[Unit]
Description=axentx TRACK B — Non-IT Business Pipeline (biz-research → biz-pitch → biz-plan)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/surrogate-1-harvest
EnvironmentFile=/etc/surrogate-coordinator.env
Environment=PYTHONUNBUFFERED=1
Environment=REPO_ROOT=/opt/surrogate-1-harvest
ExecStart=/opt/surrogate-1-harvest/.venv/bin/python /opt/surrogate-1-harvest/bin/axentx-biz-pipeline-daemon.py
Restart=always
RestartSec=30
MemoryMax=256M
TasksMax=8
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
