# systemd-dropins/axentx-trend-forecaster-daemon.service.d/memory.conf
[Service]
MemoryMax=512M
MemoryHigh=400M
